# pi-usb-automount@.service
[Service]
Type=forking
GuessMainPID=no
ExecStart=/usr/bin/pi-usb-automount %I
